# crontab
# sample crontab entry. replace the ellipses with actual path

*/1 * * * * python .../electric2go/download.py all >> .../electric2go/data/cronlog

